# files/etc/uci-defaults/94-mount-usb-docker
#!/bin/sh
# 94-mount-usb-docker
# 自动挂载 USB 硬盘，配置 Docker 数据目录到 USB

LOG_TAG="mount-usb-docker"

log() {
    echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"
}

# 检查是否有 USB 存储设备
USB_DEVICE=""
for dev in /dev/sd[a-z] /dev/sd[a-z][0-9]; do
    [ -b "$dev" ] && USB_DEVICE="$dev" && break
done

if [ -z "$USB_DEVICE" ]; then
    log "⚠️  未检测到 USB 存储设备"
    log "ℹ️  Docker 将使用 overlay 分区（约 100MB），建议插入 USB 硬盘"
    exit 0
fi

log "检测到 USB 设备: $USB_DEVICE"

# 自动挂载 USB
MOUNT_POINT="/mnt/usb_docker"
mkdir -p "$MOUNT_POINT"

# 尝试 ext4 或 NTFS
if mount -t ext4 "$USB_DEVICE" "$MOUNT_POINT" 2>/dev/null; then
    log "✅ USB 已挂载为 ext4"
elif mount -t ntfs-3g "$USB_DEVICE" "$MOUNT_POINT" 2>/dev/null; then
    log "✅ USB 已挂载为 NTFS"
else
    log "❌ USB 挂载失败"
    exit 1
fi

# 创建 Docker 数据目录
DOCKER_DATA="$MOUNT_POINT/docker"
mkdir -p "$DOCKER_DATA"
chmod 755 "$DOCKER_DATA"

# 配置 Docker 数据目录
if [ -f /etc/docker/daemon.json ]; then
    # 备份原配置
    mv /etc/docker/daemon.json /etc/docker/daemon.json.backup
fi

cat > /etc/docker/daemon.json <<EOF
{
  "data-root": "$DOCKER_DATA",
  "log-driver": "json-file",
  "log-opts": {
    "max-size": "10m",
    "max-file": "3"
  },
  "storage-driver": "overlay2",
  "iptables": true,
  "ip-forward": true
}
EOF

log "✅ Docker 数据目录配置为: $DOCKER_DATA"

# 创建其他数据目录
mkdir -p "$MOUNT_POINT/samba"  # Samba 共享
mkdir -p "$MOUNT_POINT/adguardhome"  # AdGuardHome 数据
mkdir -p "$MOUNT_POINT/downloads"    # 下载目录

log "✅ USB 数据目录结构已创建"

# 配置 Samba 默认共享目录
if [ -f /etc/samba/smb.conf.template ]; then
    sed -i "s|/mnt/sda1|$MOUNT_POINT/samba|g" /etc/samba/smb.conf.template
    log "✅ Samba 共享目录更新为 USB"
fi

# 创建标记文件
echo "USB_DOCKER_MOUNTED=1" > /etc/usb-docker-status
echo "USB_PATH=$MOUNT_POINT" >> /etc/usb-docker-status
echo "DOCKER_DATA=$DOCKER_DATA" >> /etc/usb-docker-status

log "USB Docker 配置完成"
exit 0
