# files/etc/uci-defaults/95-install-extras
#!/bin/sh
# 95-install-extras
# 后台安装大体积软件包，优化安装顺序和错误处理

LOG_TAG="install-extras"
DONE_FLAG="/etc/extras-install-done"
LOG_FILE="/var/log/extras-install.log"

# 已完成则跳过
[ -f "$DONE_FLAG" ] && exit 0

# 后台执行安装
nohup sh -c '
log() {
    echo "[$(date "+%H:%M:%S")] $1" | tee -a "$LOG_FILE" | logger -t "$LOG_TAG"
}

log "开始安装额外软件包..."

# 等待系统完全启动
sleep 30

# 等待网络就绪（最多5分钟）
log "等待网络连接..."
for i in $(seq 1 150); do
    wget -q --spider https://github.com 2>/dev/null && break
    sleep 2
done

# 更新 apk 索引
log "更新软件包索引..."
apk update 2>&1 | tee -a "$LOG_FILE" || true

# 优化安装顺序：基础依赖 -> 大体积包 -> UI 包
PACKAGES=(
    # 基础依赖
    "samba4-server"
    # Docker 相关（已在 96-install-docker 中处理）
    "adguardhome"
    # Lucky 内网穿透
    "lucky luci-app-lucky luci-i18n-lucky-zh-cn"
    # 文件管理器
    "luci-app-fileassistant"
)

for pkg_group in "${PACKAGES[@]}"; do
    log "安装: $pkg_group"
    if apk add --no-cache $pkg_group 2>&1 | tee -a "$LOG_FILE"; then
        log "✅ $pkg_group 安装成功"
    else
        log "⚠️  $pkg_group 安装失败"
        # 继续尝试下一个包
    fi
    sleep 2
done

# 刷新 LuCI 缓存
log "刷新 LuCI 界面..."
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart >/dev/null 2>&1 || true

# 创建完成标记
touch "$DONE_FLAG"
log "额外软件包安装完成"
log "安装日志: $LOG_FILE"
' > "$LOG_FILE" 2>&1 &

exit 0
