# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras
# 两阶段：① 扩容 overlay 到 U 盘并重启 ② 安装重型应用
# 依赖：编译时需要 e2fsprogs (mkfs.ext4, blkid)

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ] && exit 0
touch "$LOCK"

log "▶ install-extras 启动"

# ── 清理 video feed ──────────────────────────────────────────────────
for f in /etc/apk/repositories.d/*.list; do
    [ -f "$f" ] && sed -i '/\/video\//d' "$f"
done

# ── 等待网络 ─────────────────────────────────────────────────────────
sleep 10
NETWORK_OK=0
for i in $(seq 1 60); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    ping -c 1 -W 2 8.8.8.8  >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    sleep 3
done
if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络不可达，等下次 WAN UP 重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ 网络就绪"

printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
log "✅ DNS 注入：223.5.5.5 / 8.8.8.8"

# ════════════════════════════════════════════════════════════════════
# 辅助函数：检查是否已使用 U 盘 overlay
# 直接检查 /overlay 所在的设备是否为 /dev/sdX 或 /dev/mmcblkX
# ════════════════════════════════════════════════════════════════════
is_usb_overlay() {
    local dev
    dev=$(df /overlay 2>/dev/null | tail -1 | awk '{print $1}')
    echo "$dev" | grep -qE '^/dev/sd[a-z]|^/dev/mmcblk'
}

# 检查扩容标记（U盘上的标志 或 overlay 已生效）
if [ -f "/overlay/.overlay-active" ] || is_usb_overlay; then
    log "✅ 检测到扩容标记，进入第二阶段：安装应用"
    PHASE=2
    # 第二阶段直接使用 /overlay 路径，不需要单独定义 USB_MOUNT
else
    log "ℹ️ 未检测到扩容标记，进入第一阶段：扩容 overlay"
    PHASE=1
fi

# ════════════════════════════════════════════════════════════════════
# 第一阶段：扩容 overlay 到 U 盘
# ════════════════════════════════════════════════════════════════════
if [ "$PHASE" -eq 1 ]; then

    # 等待 U 盘设备节点出现
    USB_DEV=""
    log "⏳ 等待 U 盘设备..."
    for i in $(seq 1 15); do
        if [ -b /dev/sda1 ]; then
            USB_DEV="/dev/sda1"; break
        fi
        if [ -b /dev/sdb1 ]; then
            USB_DEV="/dev/sdb1"; break
        fi
        sleep 2
    done

    if [ -z "$USB_DEV" ]; then
        log "⚠️ 未检测到 U 盘，跳过扩容"
        log "ℹ️ 插入 U 盘后重启路由器可重新触发"
        rm -f "$LOCK"; exit 0
    fi

    log "🔍 发现设备 $USB_DEV，检查文件系统..."

    # 检查 mkfs.ext4 是否可用
    if ! command -v mkfs.ext4 >/dev/null 2>&1; then
        log "❌ mkfs.ext4 不存在，请确保编译时包含 e2fsprogs"
        log "ℹ️ 尝试安装 e2fsprogs..."
        apk add --no-cache e2fsprogs 2>&1 | logger -t "$LOG_TAG"
        if ! command -v mkfs.ext4 >/dev/null 2>&1; then
            log "❌ 仍无法找到 mkfs.ext4，放弃扩容"
            rm -f "$LOCK"; exit 1
        fi
    fi

    # 检查是否已经是 ext4，不是则格式化
    FS_TYPE=$(blkid -s TYPE -o value "$USB_DEV" 2>/dev/null)
    if [ "$FS_TYPE" != "ext4" ]; then
        log "🔧 文件系统为 ${FS_TYPE:-未知}，即将格式化为 ext4..."
        log "⚠️ U 盘数据将被清空，5 秒后开始..."
        sleep 5
        mkfs.ext4 -F -L "owrt-overlay" "$USB_DEV" 2>&1 | logger -t "$LOG_TAG"
        if [ $? -ne 0 ]; then
            log "❌ 格式化失败，放弃扩容"
            rm -f "$LOCK"; exit 1
        fi
        log "✅ 格式化完成"
    else
        log "✅ 已是 ext4，跳过格式化"
    fi

    # 挂载到临时目录
    TEMP_MOUNT="/mnt/usb-setup"
    mkdir -p "$TEMP_MOUNT"
    mount "$USB_DEV" "$TEMP_MOUNT" 2>&1 | logger -t "$LOG_TAG"
    if ! mountpoint -q "$TEMP_MOUNT"; then
        log "❌ 挂载失败，放弃扩容"
        rm -f "$LOCK"; exit 1
    fi

    # 读取 UUID
    USB_UUID=$(blkid -s UUID -o value "$USB_DEV" 2>/dev/null)
    if [ -z "$USB_UUID" ]; then
        log "❌ 无法获取 UUID，放弃扩容"
        umount "$TEMP_MOUNT"
        rm -f "$LOCK"; exit 1
    fi
    log "✅ U 盘 UUID: $USB_UUID"

    # 准备 overlay 目录结构
    mkdir -p "$TEMP_MOUNT/upper" "$TEMP_MOUNT/work"

    # 迁移现有配置到 U 盘（可选，保留原系统配置）
    if [ -d "/overlay/upper" ] && [ "$(ls -A /overlay/upper 2>/dev/null)" ]; then
        log "📦 迁移现有配置到 U 盘..."
        cp -a /overlay/upper/. "$TEMP_MOUNT/upper/" 2>/dev/null
        log "✅ 配置迁移完成"
    fi

    # 写入扩容标记到 U 盘
    touch "$TEMP_MOUNT/.overlay-active"

    # 写入 fstab（使用 UUID 挂载，避免盘符漂移）
    cat > /etc/config/fstab << FSTAB_EOF
config global
    option anon_swap '0'
    option anon_mount '0'
    option auto_swap '1'
    option auto_mount '1'
    option delay_root '8'
    option check_fs '0'

config mount
    option target '/overlay'
    option uuid '$USB_UUID'
    option fstype 'ext4'
    option options 'rw,noatime'
    option enabled '1'
    option enabled_fsck '0'
FSTAB_EOF

    uci commit fstab
    /etc/init.d/fstab enable 2>/dev/null || true

    sync
    umount "$TEMP_MOUNT"
    log "✅ fstab 已写入 UUID=$USB_UUID"
    log "🔄 5 秒后重启，重启后自动进入第二阶段..."
    rm -f "$LOCK"
    sleep 5
    reboot
    exit 0
fi

# ════════════════════════════════════════════════════════════════════
# 第二阶段：安装应用（overlay 已在 U 盘，空间充足）
# ════════════════════════════════════════════════════════════════════

# 验证 overlay 确实在外部存储
if ! is_usb_overlay; then
    log "⚠️ 未能确认 U 盘 overlay，可能是首次启动或 U 盘异常"
    log "ℹ️ 请确保 U 盘已插入，然后执行: touch /etc/install-extras-done && reboot"
    rm -f "$LOCK"
    exit 0
fi
log "✅ 确认 U 盘 overlay 模式，空间充足"

# ── apk update ───────────────────────────────────────────────────────
apk_update_ok() {
    local out n
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0

    log "⚠️ HTTPS 包数不足（${n:-0}），切换 HTTP..."
    for f in /etc/apk/repositories.d/*.list; do
        [ -f "$f" ] && sed -i \
            's|https://downloads.immortalwrt.org|http://downloads.immortalwrt.org|g' "$f"
    done
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0
    return 1
}

if ! apk_update_ok; then
    log "❌ apk update 失败，等下次重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ apk update 成功"

# ── 安装函数（修正：正确捕获 apk 退出码）────────────────────────────────
INSTALL_FAIL=0
install_pkg() {
    local out ret
    log "📦 安装: $*"
    out=$(apk add --no-cache "$@" 2>&1)
    ret=$?
    echo "$out" | logger -t "$LOG_TAG"
    if [ "$ret" -eq 0 ]; then
        log "✅ 完成: $*"
    else
        log "⚠️ 失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    fi
}

# ── 所有应用安装（overlay 在 U 盘，空间充足）────────────────────────────

# 基础工具
install_pkg smartmontools bash htop fdisk lsblk

# Samba
install_pkg samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn

# Docker（分组处理，dockerd 需要 --force-broken-world）
install_pkg docker

log "📦 强制安装: dockerd luci-app-dockerman luci-i18n-dockerman-zh-cn"
out=$(apk add --no-cache --force-broken-world \
    dockerd luci-app-dockerman luci-i18n-dockerman-zh-cn 2>&1)
ret=$?
echo "$out" | logger -t "$LOG_TAG"
if [ "$ret" -eq 0 ]; then
    log "✅ 完成: dockerd luci-app-dockerman"
else
    log "⚠️ 失败: dockerd luci-app-dockerman"
    INSTALL_FAIL=$((INSTALL_FAIL + 1))
fi

# Docker 数据目录（直接在 /overlay 下）
DOCKER_DATA="/overlay/docker-data"
mkdir -p "$DOCKER_DATA" /etc/docker
cat > /etc/docker/daemon.json << EOF
{
  "data-root": "$DOCKER_DATA",
  "log-driver": "json-file",
  "log-opts": {"max-size": "10m", "max-file": "3"}
}
EOF
log "✅ Docker data-root: $DOCKER_DATA"

# 启动 Docker
mountpoint -q /sys/fs/cgroup || mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true
/etc/init.d/dockerd enable 2>/dev/null || true
/etc/init.d/dockerd start 2>/dev/null || true

# AdGuardHome
install_pkg adguardhome luci-app-adguardhome

# AdGuardHome 数据目录（直接在 /overlay 下）
AGH_DIR="/overlay/adguardhome"
mkdir -p "$AGH_DIR"/{data,config}

# 停止服务以便配置
/etc/init.d/adguardhome stop 2>/dev/null || true

# 迁移现有配置（如果存在且目标为空）
if [ ! "$(ls -A $AGH_DIR/config 2>/dev/null)" ] && \
   [ -f "/etc/adguardhome/adguardhome.yaml" ]; then
    cp -a /etc/adguardhome/. "$AGH_DIR/config/"
    log "✅ AdGuardHome 配置已迁移到 $AGH_DIR"
fi

# 配置 UCI
uci set adguardhome.config.config_file="$AGH_DIR/config/adguardhome.yaml"
uci set adguardhome.config.work_dir="$AGH_DIR/data"
uci commit adguardhome
log "✅ AdGuardHome 配置已指向: $AGH_DIR"

# 启动 AdGuardHome
/etc/init.d/adguardhome enable 2>/dev/null || true
/etc/init.d/adguardhome start 2>/dev/null || true

# ── 刷新 LuCI ────────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

# 清理触发器
rm -f /etc/hotplug.d/iface/99-install-extras
rm -f "$LOCK"

# ── 写 DONE_FLAG ─────────────────────────────────────────────────────
if [ "$INSTALL_FAIL" -eq 0 ]; then
    touch "$DONE_FLAG"
    log "🎉 install-extras 全部完成"
else
    log "⚠️ 有 $INSTALL_FAIL 组包失败，下次 WAN UP 重试"
    mkdir -p /etc/hotplug.d/iface
    cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
    chmod +x /etc/hotplug.d/iface/99-install-extras
fi
RUNEOF

chmod +x /usr/lib/install-extras/run.sh

# ── 安装 hotplug 触发器 ──────────────────────────────────────────────
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF

chmod +x /etc/hotplug.d/iface/99-install-extras

logger -t "96-install-extras" "✅ hotplug 触发器已安装"
exit 0
