# files/etc/uci-defaults/96-install-docker
#!/bin/sh
# 96-install-docker
# 等待 USB 挂载后再安装和启动 Docker

LOG_TAG="install-docker"
DONE_FLAG="/etc/docker-install-done"
LOG_FILE="/var/log/docker-install.log"

[ -f "$DONE_FLAG" ] && exit 0

nohup sh -c '
log() {
    echo "[$(date "+%H:%M:%S")] $1" | tee -a "$LOG_FILE" | logger -t "$LOG_TAG"
}

log "Docker 安装流程开始..."

# 等待 USB 挂载配置完成（94-mount-usb-docker）
log "等待 USB 挂载配置..."
for i in $(seq 1 60); do
    if [ -f /etc/usb-docker-status ]; then
        . /etc/usb-docker-status
        if [ "$USB_DOCKER_MOUNTED" = "1" ]; then
            log "✅ USB 已挂载: $USB_PATH"
            log "✅ Docker 数据目录: $DOCKER_DATA"
            break
        fi
    fi
    sleep 5
done

# 如果没有 USB，警告但继续
if [ ! -f /etc/usb-docker-status ] || [ "$USB_DOCKER_MOUNTED" != "1" ]; then
    log "⚠️  未检测到 USB，Docker 数据将存储在 overlay（约 100MB）"
    log "⚠️  建议插入 USB 硬盘以获得完整 Docker 功能"
fi

# 更新软件包索引
log "更新软件包索引..."
apk update >/dev/null 2>&1 || true

# 安装 Docker（已在固件中编译，这里只是确保）
log "验证 Docker 安装..."
if command -v dockerd >/dev/null 2>&1; then
    log "✅ Docker 已编译进固件"
else
    log "❌ Docker 未找到，尝试 apk 安装..."
    apk add --no-cache dockerd docker docker-compose 2>&1 | tee -a "$LOG_FILE" || true
fi

# 启动 Docker
log "启动 Docker 服务..."
service dockerd enable >/dev/null 2>&1
if service dockerd start 2>&1 | tee -a "$LOG_FILE"; then
    sleep 5
    if docker info >/dev/null 2>&1; then
        log "✅ Docker 启动成功"
        log "数据目录: $(docker info 2>/dev/null | grep "Docker Root Dir" | cut -d: -f2 | tr -d " ")"
    else
        log "⚠️  Docker 未响应"
    fi
else
    log "❌ Docker 启动失败"
fi

touch "$DONE_FLAG"
log "Docker 安装完成"
' > "$LOG_FILE" 2>&1 &

exit 0
