# files/etc/uci-defaults/98-network-check
#!/bin/sh
# 98-network-check
# 验证网络服务正常工作，确保刷机后能获取 IP

LOG_TAG="network-check"

log() {
    echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"
}

log "=== 网络服务状态检查 ==="

# 1. 检查 dnsmasq (DHCP/DNS)
if pgrep -x dnsmasq >/dev/null; then
    log "✅ dnsmasq (DHCP/DNS) 正在运行"
else
    log "⚠️  dnsmasq 未运行，尝试启动..."
    if /etc/init.d/dnsmasq restart 2>/dev/null; then
        sleep 2
        if pgrep -x dnsmasq >/dev/null; then
            log "✅ dnsmasq 启动成功"
        else
            log "❌ dnsmasq 启动失败"
        fi
    else
        log "❌ dnsmasq 服务不存在"
    fi
fi

# 2. 检查防火墙
if command -v nft >/dev/null 2>&1; then
    if nft list ruleset 2>/dev/null | grep -q "accept.*established"; then
        log "✅ 防火墙规则允许已建立的连接"
    else
        log "⚠️  防火墙规则可能阻止 LAN 访问"
    fi
else
    log "⚠️  nftables 未找到，使用 iptables"
    if iptables -L FORWARD 2>/dev/null | grep -q "ACCEPT.*state.*RELATED,ESTABLISHED"; then
        log "✅ iptables 允许已建立的连接"
    fi
fi

# 3. 检查桥接接口
if command -v brctl >/dev/null 2>&1; then
    if brctl show 2>/dev/null | grep -q "br-lan"; then
        log "✅ br-lan 桥接接口已创建"
    else
        log "❌ br-lan 桥接接口未创建"
    fi
fi

# 4. 检查网络接口
log "=== 网络接口信息 ==="
if ip link show br-lan 2>/dev/null | grep -q "UP"; then
    log "✅ br-lan 接口已启动"
    IP_ADDR=$(ip addr show br-lan 2>/dev/null | grep -o "inet [0-9.]*/" | cut -d' ' -f2 | cut -d'/' -f1)
    if [ -n "$IP_ADDR" ]; then
        log "✅ br-lan IP: $IP_ADDR"
    else
        log "⚠️  br-lan 没有 IP 地址"
    fi
else
    log "❌ br-lan 接口未启动"
fi

# 5. 输出网络状态摘要
log "=== 网络状态摘要 ==="
echo "LAN IP: $(uci get network.lan.ipaddr 2>/dev/null || echo '未知')" | logger -t "$LOG_TAG"
echo "DHCP: $(uci get dhcp.lan.start 2>/dev/null || echo '未知')-$(uci get dhcp.lan.limit 2>/dev/null || echo '未知')" | logger -t "$LOG_TAG"

# 6. 创建网络就绪标记
touch /etc/network-check-done
log "网络检查完成"

exit 0
