# files/etc/uci-defaults/99-custom-ip
#!/bin/sh
# 99-custom-ip - SAFE VERSION
# 设置 LAN IP 为 192.168.2.1，验证设置成功
# 即使设置失败也能安全恢复

LOG_TAG="custom-ip"

log() {
    echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"
}

log "开始配置 LAN IP..."

# 备份当前配置
uci export network > /tmp/network-backup-$(date +%s).uci 2>/dev/null

# 设置新的 LAN IP
uci -q batch <<'EOF'
    # 设置管理 IP
    set network.lan.ipaddr='192.168.2.1'
    set network.lan.netmask='255.255.255.0'
    
    # 确保 DHCP 配置正确
    set dhcp.lan.start='100'
    set dhcp.lan.limit='150'
    set dhcp.lan.leasetime='12h'
    
    commit network
    commit dhcp
EOF

SET_RESULT=$?

# 验证设置
if [ $SET_RESULT -eq 0 ] && uci get network.lan.ipaddr 2>/dev/null | grep -q '192\.168\.2\.1'; then
    log "✅ LAN IP 成功设置为 192.168.2.1"
    log "✅ DHCP 范围: 192.168.2.100-192.168.2.250"
    
    # 重启网络服务
    /etc/init.d/network restart >/dev/null 2>&1 &
    
    # 创建标记文件
    echo "LAN IP configured: 192.168.2.1/24" > /etc/lan-ip-configured
    exit 0
else
    log "❌ LAN IP 设置失败，保持默认配置"
    
    # 恢复配置
    uci -q batch <<'EOF'
        delete network.lan.ipaddr
        delete network.lan.netmask
        commit network
EOF
    
    # 重启网络服务恢复默认
    /etc/init.d/network restart >/dev/null 2>&1 &
    
    # 创建备用标记
    echo "LAN IP: default (check with 'uci get network.lan.ipaddr')" > /etc/lan-ip-configured
    exit 0
fi
